# deploy/systemd/ai-research-celery.service
[Unit]
Description=AI Research System Celery Worker
After=network.target

[Service]
Type=forking
User=www-data
Group=www-data
EnvironmentFile=/var/www/ai_research_system/.env
WorkingDirectory=/var/www/ai_research_system
ExecStart=/var/www/ai_research_system/venv/bin/celery multi start worker1 \
          -A ai_research_system \
          --pidfile=/var/run/celery/%n.pid \
          --logfile=/var/log/celery/%n%I.log \
          --loglevel=INFO
ExecStop=/var/www/ai_research_system/venv/bin/celery multi stopwait worker1 \
         --pidfile=/var/run/celery/%n.pid
ExecReload=/var/www/ai_research_system/venv/bin/celery multi restart worker1 \
           -A ai_research_system \
           --pidfile=/var/run/celery/%n.pid \
           --logfile=/var/log/celery/%n%I.log \
           --loglevel=INFO

[Install]
WantedBy=multi-user.target