# deploy/systemd/ai-research-system.service
[Unit]
Description=AI Research System Gunicorn daemon
Requires=ai-research-system.socket
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
RuntimeDirectory=ai_research_system
WorkingDirectory=/var/www/ai_research_system
ExecStart=/var/www/ai_research_system/venv/bin/gunicorn \
          --config /var/www/ai_research_system/deploy/gunicorn.conf.py \
          ai_research_system.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Environment=DJANGO_SETTINGS_MODULE=ai_research_system.settings.production

[Install]
WantedBy=multi-user.target